# items_api/items-api.service
[Unit]
Description=Items API FastAPI Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/items_api
Environment="PATH=/path/to/items_api/venv/bin"
ExecStart=/path/to/items_api/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

